# stage-openmower/40-openmower/files/etc/systemd/system/openmower.service
[Unit]
Description=OpenMower (Docker)
Wants=network-online.target
After=docker.service network-online.target NetworkManager.service

[Service]
Type=simple
User=openmower
Group=openmower
SupplementaryGroups=docker
Restart=always
RestartSec=15s
RestartPreventExitStatus=77
TimeoutStartSec=1h
TimeoutStopSec=120s
ExecStart=/usr/local/sbin/openmower-service.sh
ExecStop=/usr/bin/docker ps -q --filter "name=openmower" | xargs -r /usr/bin/docker stop
ExecStopPost=/usr/bin/docker ps -q --filter "name=openmower" | xargs -r /usr/bin/docker rm -f

[Install]
WantedBy=multi-user.target
